# server/claude-studio-server.service
[Unit]
Description=ClaudeStudio Auth & Sync Server
After=network.target

[Service]
Type=simple
User=claude-studio
Group=claude-studio
WorkingDirectory=/opt/claude-studio-server
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

# Load config from file
EnvironmentFile=/etc/claude-studio/server.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/claude-studio
PrivateTmp=true

[Install]
WantedBy=multi-user.target
